# services/swarm-dashboard.service
[Unit]
Description=Swarm Nexus - Dashboard (Flask via Gunicorn)
After=network.target
After=swarm-collector.service
Wants=swarm-collector.service

[Service]
WorkingDirectory=/opt/swarmnexus/dashboard
Environment=FLASK_ENV=production
ExecStart=/opt/swarmnexus/dashboard/.venv/bin/gunicorn -w 3 -b 127.0.0.1:3010 app:app
Restart=always
RestartSec=5
User=www-data
Group=swarm

[Install]
WantedBy=multi-user.target
