# ops/systemd/swarm-dashboard.service.d/env.conf
[Service]
Environment=SOCIAL_TW_HANDLES=torus_network,SwarmNexus,SwarmSocials,TorusAgent
Environment=SOCIAL_EDITOR_ALLOWLIST=0xinverse
Environment=GCAL_ICS_URL=https://calendar.google.com/calendar/ical/f181e5e803b550e6715b6d319ed9eaf3a9ced61a712e391f5ea77eb6891a17fd%40group.calendar.google.com/public/basic.ics
Environment=MEDIUM_RSS=https://medium.com/feed/@Omni_Labs_
# Environment=EMISSIONS_MONTHLY=500
